# docker/front/config/supervisord_app.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
loglevel=debug
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[program:nginx]
command=/opt/cnaas/start-nginx.sh
autorestart=true
